# templates/cron/certbot-renew.cron
# ============================================
# Certbot 证书续期 Cron 任务 (自动生成)
# 每月 1 号凌晨 2 点执行
# ============================================
0 2 1 * * root @@SCRIPT_RENEW_CERT@@ >> @@CERT_RENEW_LOG_FILE@@ 2>&1
