# templates/cron/changeip.cron
# ============================================
# ChangeIP DDNS Cron 任务 (自动生成)
# 每 15 分钟执行一次
# ============================================
*/15 * * * * root @@SCRIPT_DDNS_CHANGEIP@@ >/dev/null 2>&1
